# empire-storm-strike.service
[Unit]
Description=Empire OS job: empire-storm-strike
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/empire_os
EnvironmentFile=/root/empire_os/.env
ExecStart=/root/empire_os/venv/bin/python3 /root/empire_os/storm_strike.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
