# empire-seo-loop.service
[Unit]
Description=Empire OS job: empire-seo-loop (daily 18:00)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/.hermes/scripts
ExecStart=/root/.hermes/scripts/seo_loop.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
