# empire-outreach-push.timer
[Unit]
Description=Timer for Empire OS job: empire-outreach-push (daily 07:00)

[Timer]
OnCalendar=07:00
Persistent=true
Unit=empire-outreach-push.service

[Install]
WantedBy=timers.target
