# empire-settle-funnel.service
[Unit]
Description=Empire OS job: empire-settle-funnel (daily 19:00)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/.hermes/scripts
ExecStart=/root/.hermes/scripts/settle_funnel.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
