# empire-supervisor-alert.timer
[Unit]
Description=Timer for Empire OS job: empire-supervisor-alert (every 15 min)

[Timer]
OnCalendar=*-*:00/15
Persistent=true
Unit=empire-supervisor-alert.service

[Install]
WantedBy=timers.target
